# systemd/socket-bridge@.service
[Unit]
Description=Unix Socket Bridge - %i
After=network.target

[Service]
Type=simple
User=socket-bridge
Group=socket-bridge
ExecStart=/usr/local/bin/socket-server /etc/socket-bridge/%i.json
Restart=always
RestartSec=5

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

# Runtime directory
RuntimeDirectory=socket-bridge
RuntimeDirectoryMode=0755

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
