# systemd/socket-bridge@.service
[Unit]
Description=Unix Socket Bridge - %i
After=network.target

[Service]
Type=simple
User=your-username
Group=your-username
# Add session environment for DBUS access (adjust UID as needed)
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
ExecStart=/usr/local/bin/unix-socket-server /etc/socket-bridge/%i.json
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=10

# Security settings
NoNewPrivileges=true

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
